# university_of_washington_cse_mushr/mushr_ws/devel/etc/catkin/profile.d/cse478.bash
/home/robotics/mushr_ws/devel/.private/cse478/etc/catkin/profile.d/cse478.bash